# plugins/community/root/__deleted/misaki_yooko/workspace/config/init.conf

[program:misaki_yooko-c]
user=root
directory=/home/karim/Workspace/github.com/vcar/vcar/plugins/community/root/misaki_yooko/services
command=/home/karim/Workspace/github.com/vcar/vcar/plugins/community/root/misaki_yooko/workspace/.env/bin/gunicorn --bind unix:app.sock --reload wsgi
; -e SCRIPT_NAME=/plugins/community/root/misaki_yooko/ 
autostart=true
autorestart=true

stderr_logfile=/usr/local/var/log/vcar/misaki_yooko/misaki_yooko.err.log
stderr_logfile_maxbytes=2MB
stdout_logfile=/usr/local/var/log/vcar/misaki_yooko/misaki_yooko.out.log
stdout_logfile_maxbytes=2MB