# plugins/core/p1/workspace/config/init.conf
[program:p1@plugin.core]

autostart=true
autorestart=true


directory=/usr/src/app/plugins/core/p1/services
command=gunicorn --workers 3 --bind unix:app.sock --reload wsgi

stderr_logfile=/var/log/supervisor/vcar/p1.err.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/vcar/p1.out.log
stdout_logfile_maxbytes=10MB