# plugins/community/root/gPI/workspace/config/init.conf
[program:gPI-c]
user=root
directory=/Users/salaheddine/Desktop/vcar_repository/vcar/plugins/community/root/gPI/services
command=/Users/salaheddine/Desktop/vcar_repository/vcar/plugins/community/root/gPI/workspace/.env/bin/gunicorn --bind unix:app.sock --reload wsgi
; -e SCRIPT_NAME=/plugins/community/root/gPI/ 
autostart=true
autorestart=true

stderr_logfile=/usr/local/var/log/vcar/gPI/gPI.err.log
stderr_logfile_maxbytes=2MB
stdout_logfile=/usr/local/var/log/vcar/gPI/gPI.out.log
stdout_logfile_maxbytes=2MB