# plugins/community/root/p1/workspace/config/init.conf
[group:yooko]
programs=xx

[program:xx]
user=root
directory=/home/karim/Workspace/github.com/vcar/vcar/plugins/community/root/p1/services
command=/home/karim/Workspace/github.com/vcar/vcar/plugins/community/root/p1/workspace/.env/bin/gunicorn --bind unix:app.sock --reload wsgi
; -e SCRIPT_NAME=/plugins/community/root/p1/
autostart=true
autorestart=true

stderr_logfile=/var/log/vcar/p1/p1.err.log
stderr_logfile_maxbytes=2MB
stdout_logfile=/var/log/vcar/p1/p1.out.log
stdout_logfile_maxbytes=2MB